# deploy/systemd/punchgate.service
[Unit]
Description=Punchgate background daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
EnvironmentFile=%h/.config/punchgate/punchgate.env
ExecStart=/usr/local/bin/punchgate
Restart=always
RestartSec=2
KillSignal=SIGINT
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
