# deploy/launchd/com.punchgate.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.punchgate.daemon</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/punchgate</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PUNCHGATE_IDENTITY</key>
    <string>/Users/you/.config/punchgate/identity.key</string>
    <key>PUNCHGATE_LISTEN</key>
    <string>/ip4/0.0.0.0/udp/0/quic-v1</string>
    <key>RUST_LOG</key>
    <string>info</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/you/Library/Logs/punchgate.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/you/Library/Logs/punchgate.err.log</string>
</dict>
</plist>
